# lede-ecr-app/files/uci-defaults/51-fix-ppp
#!/bin/sh

cp /home/ftp_user/ppp.sh /lib/netifd/proto/ppp.sh