# lede-ecr-app/files/hotplug.d/usb/30_usb-hid
logger -t HOTPLUG "Always $PRODUCT $ACTION"
if [ "${ACTION}" = "add" ]; then
        logger -t HOTPLUG "ADD new: $PRODUCT"
        /usr/bin/node /www/ecr-app/hid.js $PRODUCT
fi

