# miniupnpd/files/uci-defaults/51-friendly-name
#!/bin/sh

echo friendly_name=Omega-$(cat /sys/class/net/eth0/address|awk -F ":" '{print $5""$6}'| tr a-z A-Z) >> /usr/share/miniupnpd/miniupnpd.manufacture.conf